# deploy/weather-display.service
[Unit]
Description=Raspberry Pi E‑Ink Weather Display (venv)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/weather-display
ExecStart=/home/pi/weather-display/.venv/bin/python /home/pi/weather-display/main.py \
          --config /home/pi/weather-display/config.yaml
Restart=on-failure
RestartSec=10
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target